# modules/parity/files/systemd/parity.service
[Unit]
Description=Parity Daemon
 
[Service]
User=parity
ExecStart=/usr/bin/parity
 
[Install]
WantedBy=default.target


